# systemd/.config/systemd/user/isync.service
[Unit]
Description=Mailbox Sync 

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -VXa
ExecStartPost=/usr/bin/notmuch new
