# systemd/.config/systemd/user/isync.timer
[Unit]
Description=Mailbox Sync Timer

[Timer]
OnBootSec=1m
OnUnitActiveSec=2m
Unit=isync.service

[Install]
WantedBy=timers.target
