# systemd/.config/systemd/user/mailcheck.service
[Unit]
Description=Mail Checker Daemon 

[Service]
ExecStart=python3 -u /home/ak/bin/mailchecker
Restart=always

[Install]
WantedBy=default.target
