# systemd/.config/systemd/user/new_mails.service
[Unit]
Description=Mail Check Daemon 

[Service]
ExecStart=/home/ak/bin/new_mails
Restart=always

[Install]
WantedBy=default.target
