# systemd/.config/systemd/user/wallpaper.timer
[Unit]
Description=Random wallpaper with feh

[Timer]
OnUnitActiveSec=5m
Unit=wallpaper.service

[Install]
WantedBy=default.target
